# config/gunicorn.socket
[Unit]
Description=gunicorn socket

[Socket]
ListenStream=/run/gunicorn.sock
SocketUser=cheechuen

[Install]
WantedBy=sockets.target 
